# ubuntu_16/crontab
@reboot sleep 30 && /home/talisman/ligar_servidor.sh > /home/talisman/log_boot.txt 2>&1